# pypiserver/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}

[program:pypiserver]
command=/srv/pypiserver/bin/gunicorn --paste /srv/pypiserver/paste.ini
stdout_logfile=/var/log/pypiserver/access.log
stderr_logfile=/var/log/pypiserver/error.log
user=pypiserver
autostart=true
autorestart=true